# autoupdate/cron.d/autoupdate_cimon
# Copyright (C) Schweizerische Bundesbahnen SBB, 2016
# cimon, SBB, FSe 2016
# daily at 6:42 check for updates of the cimon controller and setup
42 6 * * * pi /opt/cimon/autoupdate/autoupdate.sh > /var/log/cimon/autoupdate_cimon.log 2>&1
