# controller/cron.d/cimon
# Copyright (C) Schweizerische Bundesbahnen SBB, 2016
# cimon, SBB, FSe 2016
# all 5 minutes start cimon controller if not running
*/5 * * * * pi /opt/cimon/watchdog/watchdog.sh > /dev/null 2>&1
# reboot cimon box every sunday at 3:42 at night
42 3 * * 6 root reboot
